# tamper.timer
[Unit]
Description=Run tamper check daily

[Timer]
OnUnitActiveSec=1day
 
[Install]
WantedBy=timers.target
